# debian/preinst
#!/bin/sh
set -e

case "$1" in
    install|upgrade)
        # Stop the service if it is running to ensure a clean binary replacement
        if [ -d /run/systemd/system ]; then
            # Stop by primary service name
            if systemctl is-active --quiet drove.gateway.service 2>/dev/null; then
                systemctl stop drove.gateway.service || true
            fi
            # Stop by alias as fallback
            if systemctl is-active --quiet nixy.service 2>/dev/null; then
                systemctl stop nixy.service || true
            fi
        fi
        
        # Migration: Check if old nixy.service was enabled and save state for postinst
        if [ "$1" = "upgrade" ]; then
            mkdir -p /var/lib/drove-gateway
            if systemctl is-enabled nixy.service 2>/dev/null; then
                touch /var/lib/drove-gateway/.nixy-was-enabled
            fi
        fi
    ;;
    abort-upgrade)
    ;;
    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
